# files/systemd/luna-next.service
[Unit]
Description=The Luna Next UI

[Service]
Type=notify
ExecStart=/usr/bin/luna-next --systemd
Restart=always

[Install]
WantedBy=multi-user.target
